# deploy/supervisor_conf.d/webapp.conf
[program:streamsavvy_dataprocessing]

command = gunicorn streamsavvy_dataprocessing.wsgi --log-file -